# conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
autorestart = true

[program:php5-fpm]
command = /usr/sbin/php5-fpm
user = root
autostart = true
autorestart = true

[program:memcached]
command = /usr/bin/memcached -u memcache -v
user = root
autostart = true
autorestart = true
